# Services/camcontroller-web.service
# /etc/systemd/system/camcontroller-web.service
[Unit]
Description=CamController Web Interface - Gunicorn-instans for att servera Flask kamerainstallningar
After=network-online.target
Wants=network-online.target
[Service]
User=pi
Group=pi
WorkingDirectory=/home/pi/PyRpiCamController/WebGui
# If installed via apt:
ExecStart=/usr/bin/gunicorn -w 2 -b 0.0.0.0:80 web_app:app
# Allow non-root service user to bind privileged ports like 80.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
Restart=always
RestartSec=10
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/pi/PyRpiCamController"
[Install]
WantedBy=multi-user.target
